# ansible/roles/raspberry_app_dyndns/templates/dyndns.timer.j2
[Unit]
Description={{ item.domain }} DynDs update timer

[Timer]
OnCalendar=daily
OnBootSec=10min
Persistent=false

[Install]
WantedBy=timers.target
