# roles/raspberry_app_crm/templates/crm-cron.timer.j2
[Unit]
Description=CRM Cron

[Timer]
OnCalendar=hourly
OnBootSec=10min
Persistent=false

[Install]
WantedBy=timers.target
