# ansible/roles/ZZ-borgmatic/templates/borgmatic.timer.j2
[Unit]
Description=Borgmatic backup

[Timer]
OnCalendar=daily
OnBootSec=10min
Persistent=false

[Install]
WantedBy=timers.target
